# auto/ba/ba-banner.timer
[Unit]
Description=Timer for ba-banner.service

[Timer]
OnUnitActiveSec=3h
OnBootSec=1min

[Install]
WantedBy=timers.target

